# Virtual/Machine/ISO-maker/Files/For-Target/sys_root_dir/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --noissue --autologin econ-ark %I $TERM
Type=idle
